# deploy/peerweights.service
[Unit]
Description=PeerWeights API Server
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/peerweights/server
ExecStart=/usr/bin/node --import tsx src/index.ts
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/peerweights/server/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/peerweights

[Install]
WantedBy=multi-user.target
